# deploy/vehicle-scraper.service
[Unit]
Description=Vehicle Market Price Scraper
After=network-online.target postgresql.service

[Service]
Type=oneshot
User=scraper
WorkingDirectory=/opt/mobile-scraper
EnvironmentFile=/opt/mobile-scraper/.env
ExecStart=/usr/bin/xvfb-run -a --server-args="-screen 0 1920x1080x24" /opt/mobile-scraper/venv/bin/python main.py
ExecStartPost=/opt/mobile-scraper/venv/bin/python report.py
TimeoutStartSec=600
StandardOutput=journal
StandardError=journal
