# deploy/vehicle-report.service
[Unit]
Description=Vehicle Price Tracker Web Server
After=network-online.target postgresql.service

[Service]
Type=simple
User=scraper
WorkingDirectory=/opt/mobile-scraper
ExecStart=/opt/mobile-scraper/venv/bin/python server.py
Restart=always
RestartSec=5
Environment=PORT=8080

[Install]
WantedBy=multi-user.target
